# arch/battery-notification.service
[Unit]
Description=Battery Notification Service
[Service]
ExecStart=/path/to/battery_full_notification.sh
